# tools/depends/target/dbus/raspberry-pi2-release/bus/systemd-user/dbus.service
[Unit]
Description=D-Bus User Message Bus
Documentation=man:dbus-daemon(1)
Requires=dbus.socket

[Service]
ExecStart=/root/kodi-rpi/raspberry-pi2-release/bin/dbus-daemon --session --address=systemd: --nofork --nopidfile --systemd-activation --syslog-only
ExecReload=/root/kodi-rpi/raspberry-pi2-release/bin/dbus-send --print-reply --session --type=method_call --dest=org.freedesktop.DBus / org.freedesktop.DBus.ReloadConfig
