# tools/depends/target/dbus/raspberry-pi2-release/bus/dbus.socket
[Unit]
Description=D-Bus System Message Bus Socket

[Socket]
ListenStream=/root/kodi-rpi/raspberry-pi2-release/var/run/dbus/system_bus_socket
